# elephantbroker-hitl.service
[Unit]
Description=ElephantBroker HITL Middleware
After=elephantbroker.service

[Service]
Type=simple
User=dbadmin
EnvironmentFile=/etc/elephantbroker/hitl.env
ExecStart=/opt/elephant-broker/venv/bin/python -m hitl_middleware
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
