# deploy/systemd/elephantbroker.service
# =============================================================================
# ElephantBroker Cognitive Runtime — systemd service unit
# =============================================================================
#
# Installation:
#   sudo cp deploy/systemd/elephantbroker.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable elephantbroker
#   sudo systemctl start elephantbroker
#
# Verification:
#   systemctl status elephantbroker
#   journalctl -u elephantbroker -f
#   curl http://localhost:8420/health/        # note trailing slash
#
# This unit assumes the canonical install layout produced by deploy/install.sh:
#
#   /opt/elephantbroker/        elephantbroker:elephantbroker  755   venv lives here
#   /etc/elephantbroker/        elephantbroker:elephantbroker  750   config dir
#   /etc/elephantbroker/default.yaml   elephantbroker:elephantbroker  640
#   /etc/elephantbroker/env            root:elephantbroker            640   secrets
#   /var/lib/elephantbroker/    elephantbroker:elephantbroker  750   runtime data
#
# The `elephantbroker` system user is created by deploy/install.sh with
# `useradd --system --shell /usr/sbin/nologin --home-dir /var/lib/elephantbroker`.
# It exists ONLY to run this service — no interactive login is possible.
# =============================================================================

[Unit]
Description=ElephantBroker Cognitive Runtime
Documentation=https://github.com/elephant-broker/elephant-broker
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple

# --- Identity ---
# Dedicated system user. Created by deploy/install.sh.
User=elephantbroker
Group=elephantbroker

# --- Filesystem ---
WorkingDirectory=/var/lib/elephantbroker
EnvironmentFile=/etc/elephantbroker/env
ExecStart=/opt/elephantbroker/.venv/bin/elephantbroker serve \
    --config /etc/elephantbroker/default.yaml \
    --host 0.0.0.0 \
    --port 8420

# --- Restart policy ---
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

# --- Logging ---
StandardOutput=journal
StandardError=journal
SyslogIdentifier=elephantbroker

# --- Filesystem hardening ---
# Most of the filesystem is invisible or read-only. Only the runtime's data
# directory and the venv (where Cognee writes its internal state) are
# writable. /etc/elephantbroker remains read-only as expected.
ProtectSystem=strict
ReadWritePaths=/var/lib/elephantbroker /opt/elephantbroker
ProtectHome=true
PrivateTmp=true
PrivateDevices=true

# --- Privilege hardening ---
NoNewPrivileges=true
RestrictSUIDSGID=true
LockPersonality=true

# --- Kernel & namespace hardening ---
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# --- File creation mask ---
# New files default to 640 (owner rw, group r, other none) instead of 644.
UMask=0027

[Install]
WantedBy=multi-user.target
